# .config/systemd/user/wallp-switch.timer
[Unit]
Description=Wallpaper rotation timer

[Timer]
OnCalendar=*:0/15
Persistent=true
