# .config/yadm/alt/.config/systemd##os.Linux/user/blueman.service
[Unit]
Description=GTK+ Bluetooth Manager
Documentation=https://github.com/blueman-project/blueman
BindsTo=sway-session.target

[Service]
Type=simple
ExecStart=/usr/bin/blueman-applet

[Install]
WantedBy=sway-session.target
